# media-video/mediamtx/files/mediamtx.service
[Unit]
Wants=network.target
[Service]
ExecStart=/usr/bin/mediamtx
[Install]
WantedBy=multi-user.target