# genpack/systemimg/files/systemimg-init.service
[Unit]
Description=Initialize system.img based system

[Service]
Type=oneshot
ExecStart=/usr/lib/systemimg/init

[Install]
WantedBy=sysinit.target
